# whirlcrypt.service
[Unit]
Description=Whirlcrypt Secure File Sharing Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=andrei
Group=andrei
WorkingDirectory=/home/andrei/whirlcrypt
EnvironmentFile=/home/andrei/whirlcrypt/backend/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=whirlcrypt

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
