# whirlcrypt.service
[Unit]
Description=Whirlcrypt Secure File Sharing Service
After=network.target

[Service]
Type=simple
User=andrei
WorkingDirectory=/home/andrei/whirlcrypt/backend
EnvironmentFile=/home/andrei/whirlcrypt/backend/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=whirlcrypt

[Install]
WantedBy=multi-user.target
